# package/base-files/files/etc/uci-defaults/14_disable-dropbear
#!/bin/sh

/etc/init.d/dropbear enabled && {
	/etc/init.d/dropbear stop
	/etc/init.d/dropbear disable
}

exit 0
